# Linux/tinyvpn.service
[Unit]
Description=TinyVPN Service
After=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=bash /etc/tinyvpn/server_up.sh
ExecStop=bash /etc/tinyvpn/server_down.sh

[Install]
WantedBy=multi-user.target
